# gopanel.service
[Unit]
Description=goPanel Web Analytics

[Service]
User=root
ExecStart=/usr/local/bin/gopanel

[Install]
WantedBy=multi-user.target
